# services/kiosk.service
[Unit]
Description=Chromium kiosk
Wants=multi-user.target
After=multi-user.target

[Install]
WantedBy=multi-user.target

[Service]
Environment=DISPLAY=:0
Type=simple
ExecStart=/usr/bin/xinit ./resources/kiosk.sh
Restart=always
User=pi
Group=pi
